# src/etc/profile.d/editor.sh
export EDITOR=/usr/bin/mcedit
export VISUAL=/usr/bin/mcedit

